# config/systemd/user/chromium@.service
.include X.in

[Unit]
Description=Chromium

[Service]
ExecStart=/usr/bin/chromium --no-startup-window
Restart=always
