# config/systemd/user/compiz.service
[Unit]
Description=Compiz Window Manager

[Service]
ExecStart=/usr/bin/compiz ccp

[Install]
WantedBy=default.target
