# config/systemd/user/emerald.service
[Unit]
Description=Emerald Window Decorator

[Service]
ExecStart=/usr/bin/emerald

[Install]
WantedBy=default.target
