# config/systemd/user/enlightenment.service
[Unit]
Description=Enlightenment Window Manager

[Service]
ExecStart=/usr/bin/enlightenment_start

[Install]
WantedBy=default.target
