# config/systemd/user/herbstluftwm@.service
.include X.in

[Unit]
Description=Herbstluft Window Manager

[Service]
ExecStart=/usr/bin/herbstluftwm
