# config/systemd/user/kdeinit4@.service
.include X.in

[Unit]
Description=KDE Init

[Service]
Type=forking
ExecStart=/usr/bin/kdeinit4
