# config/systemd/user/kdeinit4.service
[Unit]
Description=KDE Init

[Service]
Type=forking
ExecStart=/usr/bin/kdeinit4

[Install]
WantedBy=default.target
