# config/systemd/user/mpd.service
.include /usr/lib/systemd/user/mpd.service

[Unit]
Wants=pulseaudio.service
