# config/systemd/user/nitrogen@.service
.include X.in

[Unit]
Description=Nitrogen Wallpaper

[Service]
Type=oneshot
ExecStart=/usr/bin/nitrogen --restore
