# config/systemd/user/numlockx@.service
.include X.in

[Unit]
Description=Keyboard NumLock

[Service]
Type=oneshot
ExecStart=/usr/bin/numlockx
