# config/systemd/user/pulseaudio.service
[Unit]
Description=PulseAudio

[Service]
ExecStart=/usr/bin/pulseaudio
Restart=always

[Install]
WantedBy=default.target
