# config/systemd/user/subtle.service
[Unit]
Description=Subtle Window Manager

[Service]
ExecStart=/usr/bin/subtle

[Install]
WantedBy=default.target
