# config/systemd/user/tint2.service
[Unit]
Description=tint2 panel/taskbar

[Service]
ExecStart=/usr/bin/tint2

[Install]
WantedBy=default.target
