# config/systemd/user/tint2@.service
.include X.in

[Unit]
Description=tint2 panel/taskbar

[Service]
ExecStart=/usr/bin/tint2
