# config/systemd/user/X.service
[Unit]
Description=X

[Service]
Type=notify
ExecStart=/usr/bin/xinit -- vt${XDG_VTNR}
ExecStopPost=/usr/bin/systemctl --user unset-environment DISPLAY
SuccessExitStatus=1

[Install]
WantedBy=graphical.target
