# config/systemd/user/X@.service
[Unit]
Description=X

[Service]
Type=forking
ExecStart=/bin/test -e /run/user/%U/X11-display
KillSignal=SIGCONT
PIDFile=/run/user/%U/xsession.pid
