# config/systemd/user/xbindkeys@.service
.include X.in

[Unit]
Description=Xbindkeys

[Service]
ExecStart=/usr/bin/xbindkeys -n
KillMode=process
