# config/systemd/user/xbindkeys.service
[Unit]
Description=Xbindkeys

[Service]
ExecStart=/usr/bin/xbindkeys
Type=forking
KillMode=process

[Install]
WantedBy=default.target
