# config/systemd/user/TODO/i3.service
[Unit]
Description=i3 Window Manager

[Service]
ExecStart=/usr/bin/i3

[Install]
WantedBy=default.target
