# config/systemd/user/TODO/razor-qt.service
[Unit]
Description=Razor-qt Desktop Environmnent

[Service]
ExecStart=/usr/bin/startrazor

[Install]
WantedBy=default.target
